# container/rootfs/services.d/apache2/run
#!/usr/bin/with-contenv sh

exec apachectl -DFOREGROUND
