# container/rootfs/services.d/apache2/run
#!/bin/sh

apachectl -D FOREGROUND
echo "apache is running"